# ansible/roles/telegraf/files/supervisord.conf
[program:telegraf]
command=/usr/bin/telegraf --config /etc/telegraf/telegraf.conf --config-directory /etc/telegraf/telegraf.d
directory=/
autostart=true
autorestart=true
startretries=9
stderr_logfile=/dev/null
stdout_logfile=/dev/null
user=telegraf
group=telegraf

