# ansible/roles/hdfs/templates/supervisord.conf
{% if "datanode" in group_names %}
[program:datanode]
command=/opt/hadoop-2.6.5/bin/hdfs datanode
directory=/opt/hadoop-2.6.5/bin
autostart=true
autorestart=true
startretries=9
redirect_stderr=true
stdout_logfile=/var/log/datanode.log
user=root
group=root
{% endif %}
{% if "namenode" in group_names %}
[program:namenode]
command=/opt/hadoop-2.6.5/bin/hdfs namenode
directory=/opt/hadoop-2.6.5/bin
autostart=true
autorestart=true
startretries=9
redirect_stderr=true
stdout_logfile=/var/log/namenode.log
user=root
group=root
{% endif %}

